# services/homepage-docker-compose.service
[Unit]
Description=Homepage docker compose service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "docker compose -f /opt/hsrv/docker-compose-homepage.yml up --detach"
ExecStop=/bin/bash -c "docker compose -f /opt/hsrv/docker-compose-homepage.yml stop"

[Install]
WantedBy=multi-user.target
