# services/portainer-docker-compose.service
[Unit]
Description=Portainer docker compose service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c "docker compose -f /opt/hsrv/docker-compose-portainer.yml up --detach"
ExecStop=/bin/bash -c "docker compose -f /opt/hsrv/docker-compose-portainer.yml stop"

[Install]
WantedBy=multi-user.target
